# deploy/systemd/sity-frontend.service
[Unit]
Description=Sity Frontend
After=network.target

[Service]
Type=simple
User=alex
WorkingDirectory=/home/alex/projects/sity/frontend
ExecStart=/usr/bin/npm run dev -- --host 0.0.0.0
Restart=on-failure
RestartSec=3
Environment=VITE_SITY_API_BASE=http://192.168.1.133:8000

[Install]
WantedBy=multi-user.target
