# deploy/audio/hdmi-audio-forward.service
[Unit]
Description=HDMI audio forward — pipeline continuo PCM→IEC958
After=pipewire.service wireplumber.service
Requires=pipewire.service wireplumber.service

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'until aplay -l 2>/dev/null | grep -q Loopback; do sleep 1; done'
ExecStart=/bin/sh -c 'arecord -D hw:Loopback,1,0 -f S16_LE -r 48000 -c 2 | \
    python3 /home/alex/.local/bin/pcm2iec958.py | \
    aplay -D hw:vc4hdmi0,0 -f IEC958_SUBFRAME_LE -r 48000 -c 2 -t raw'
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
